# modules/share-extra/airootfs.any/etc/profile.d/disable-beep.sh
#!/bin/sh
#
# silencesuzuka
# Known As: @admin
# Email  : admin@noreply
#
# (c) 1998-2140 team-silencesuzuka
#
# disable-beep.sh
#
# Script to disable PC beep sound
#

setterm -blength 0
